# deployment/tempoftp.service
[Unit]
Description=FastAPI TempoFTP Service
After=network.target

[Service]
User=lanotadm
Group=lanotadm
WorkingDirectory=/opt/tempoftp
Environment="PATH=/opt/tempoftp/venv/bin"
EnvironmentFile=/opt/tempoftp/.env
ExecStart=/opt/tempoftp/venv/bin/uvicorn main:app --workers 4 --host 127.0.0.1 --port 8000 --proxy-headers
Restart=always

[Install]
WantedBy=multi-user.target
